# network/vpn-split-tunnel-refresh.service
[Unit]
Description=Refresh VPN split tunnel bypass routes

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vpn-split-tunnel
